# services_check/docker-restart.timer
[Unit]
Description=Docker Restart Timer

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target